# etc/profile.d/nvm-proxy.sh
###
# node|npm
alias node="/usr/local/bin/node-nvm-proxy"
alias yarn="/usr/local/bin/yarn-nvm-proxy"
